# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=/usr/sbin/openvpn --cd /etc/openvpn --config RancherVPNClient.ovpn
